# deploy/kotify.service
[Unit]
Description=kotify — Korean broadcast notification system powered by NCP
Documentation=https://github.com/RAVNUS-INC/kotify
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kotify
Group=kotify
WorkingDirectory=/opt/kotify
ExecStartPre=/opt/kotify/.venv/bin/alembic -c /opt/kotify/alembic.ini upgrade head
# --proxy-headers + --forwarded-allow-ips: NPM 같은 리버스 프록시의
# X-Forwarded-Proto/For 헤더를 신뢰. 이게 없으면 SessionMiddleware가 https를
# 인식하지 못해 Secure 쿠키를 못 만들고, 결과적으로 CSRF 검증이 모두 403으로
# 차단된다. forwarded-allow-ips=* 는 모든 출처를 허용하는데, NPM 앞단에서 외부
# 인바운드를 차단하는 환경에서는 안전.
ExecStart=/opt/kotify/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080 --workers 1 --proxy-headers --forwarded-allow-ips=*

# 재시작 정책
Restart=on-failure
RestartSec=5

# 로그 출력
StandardOutput=append:/var/log/kotify/stdout.log
StandardError=append:/var/log/kotify/stderr.log

# 환경변수
# DB 경로, master.key 경로는 app/config.py 기본값(/var/lib/kotify/) 사용
# .env 파일 없음 — 모든 설정은 DB + master.key로 관리
Environment=PYTHONUNBUFFERED=1

# ── 보안 강화 (Systemd Hardening) ────────────────────────────────────────────
# NoNewPrivileges=false: 웹 UI 원클릭 업데이트 기능(/admin → 업데이트)이
# sudo로 /opt/kotify/deploy/kotify-update.sh 를 호출하는데, true일 경우 커널
# 레벨에서 setuid(sudo) 승격이 차단되어 동작 불가. 승격 범위는 kotify-sudoers
# 로 해당 스크립트 1개에만 한정되어 있으므로 실효 공격면은 좁음.
NoNewPrivileges=false
ProtectSystem=strict
# /opt/kotify 포함: 웹 UI 원클릭 업데이트(git pull + pip install -e)가 해당
# 경로에 쓰기를 요구. 없으면 systemd 네임스페이스 안에서 read-only라 스크립트가
# stderr 없이 rc=255로 실패한다.
ReadWritePaths=/var/lib/kotify /var/log/kotify /opt/kotify
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
