# deploy/kotify.service
[Unit]
Description=kotify — Korean broadcast notification system (U+ msghub)
Documentation=https://github.com/RAVNUS-INC/kotify
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kotify
Group=kotify
WorkingDirectory=/opt/kotify
ExecStartPre=/opt/kotify/.venv/bin/alembic -c /opt/kotify/alembic.ini upgrade head
# --host 127.0.0.1: FastAPI 는 내부 전용. 외부 노출은 Next.js(kotify-web.service)
# 가 담당하고, /api/* 는 next.config.mjs 의 rewrites() 로 내부 프록시된다.
# 0.0.0.0 바인딩은 같은 VLAN 의 다른 호스트가 API 를 직접 호출할 수 있어 위험.
#
# --proxy-headers + --forwarded-allow-ips=127.0.0.1: 같은 CT 안에 있는 Next.js
# 에서만 X-Forwarded-Proto/For 헤더를 신뢰한다. `*` 로 두면 위조된 헤더가
# 받아들여져 Secure 쿠키/감사 IP 가 스푸핑될 수 있다.
ExecStart=/opt/kotify/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8080 --workers 1 --proxy-headers --forwarded-allow-ips=127.0.0.1

# 재시작 정책
Restart=on-failure
RestartSec=5

# 로그 출력
StandardOutput=append:/var/log/kotify/stdout.log
StandardError=append:/var/log/kotify/stderr.log

# 환경변수
# DB 경로, master.key 경로는 app/config.py 기본값(/var/lib/kotify/) 사용
# .env 파일 없음 — 모든 설정은 DB + master.key로 관리
Environment=PYTHONUNBUFFERED=1

# ── 보안 강화 (Systemd Hardening) ────────────────────────────────────────────
# NoNewPrivileges=false: 웹 UI 원클릭 업데이트 기능(/admin → 업데이트)이
# sudo로 /opt/kotify/deploy/kotify-update.sh 를 호출하는데, true일 경우 커널
# 레벨에서 setuid(sudo) 승격이 차단되어 동작 불가. 승격 범위는 kotify-sudoers
# 로 해당 스크립트 1개에만 한정되어 있으므로 실효 공격면은 좁음.
NoNewPrivileges=false
ProtectSystem=strict
# /opt/kotify 포함: 웹 UI 원클릭 업데이트(git pull + pip install -e)가 해당
# 경로에 쓰기를 요구. 없으면 systemd 네임스페이스 안에서 read-only라 스크립트가
# stderr 없이 rc=255로 실패한다.
ReadWritePaths=/var/lib/kotify /var/log/kotify /opt/kotify
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
